# autodeploy/data-api.conf
# /etc/init/data-api.conf
# Upstart configuration for running a Node.js application.

description "Data API for mobile Meetin.gs"
author      "Tuomas Starck <tuomas@meetin.gs>"

start on    (local-filesystems and net-device-up IFACE!=lo)
stop  on    shutdown

respawn
respawn limit 5 5

setuid nobody
setgid nogroup

script
  set -e
  FIFO=`mktemp --dry-run`
  mkfifo $FIFO
  (logger -t api -p local7.info <$FIFO &)
  exec 2>$FIFO
  exec >$FIFO
  rm $FIFO

  . /etc/autodeploy.conf

  exec /usr/local/bin/meetings-data-api
end script
